# nginx/etc/systemd/system/nginx.service.d/override.conf
[Service]
LimitNOFILE=1000000
CPUQuota=33%
